# home/.config/systemd/user/break.service
[Unit]
Description=Break Service

[Service]
Type=simple
ExecStart=/usr/bin/notify-send 'Time for a Break'